# systemd/.config/systemd/user/dms.service.d/locale.conf
[Service]
Environment=LANG=zh_CN.UTF-8
Environment=LANGUAGE=zh_CN.UTF-8
